# ansible/roles/docker/templates/etc/systemd/system/docker.service.d/http-proxy.conf.j2
[Service]
Environment="HTTP_PROXY=http://10.240.68.68:3128"